# build/linux/autostart/copyd.desktop
[Desktop Entry]
Type=Application
Name=copyd (background)
Comment=Start the copyd clipboard watcher hidden in the background
Exec=copyd --hidden
Icon=copyd
Terminal=false
Categories=Utility;
X-GNOME-Autostart-enabled=true
NoDisplay=true
